# debian/install
sudo apt-get install -y fish  qutebrowser lxappearance btop nemo nitrogen arandr
sudo apt-get install -y conky picom zathura zathura-pdf-poppler openssh-server irqbalance 
sudo apt-get install -y zip unzip 7zip vlc qbittorrent
sudo apt-get install -y blender audacity godot3 gimp
sudo apt install openssh-server
sudo apt-get install dunst
sudo apt-get install plymouth plymouth-themes



sudo apt-get remove -y firefox-esr




sudo systemctl enable --now irqbalance

sudo systemctl enable ssh
sudo systemctl start ssh


mkdir -p ~/tmp
cd ~/tmp
git clone https://github.com/vim/vim.git
cd vim
./configure --with-features=huge --enable-python3interp
make
sudo make install